# 02_bash/02_bash/sed/inittab
id:2:initdefault:
id:4:initdefault:
id:4:initdefaultt:
id:4:initdefaultt:


id:4:initdefaultt:
